# keepalived/keepalived.ini
[program:keepalived]
command=/usr/sbin/keepalived -n -l -D -f /etc/keepalived/keepalived.conf --dont-fork --log-console
autostart=true
autorestart=true
stderr_logfile=/var/log/keepalived.err.log
stdout_logfile=/var/log/keepalived.out.log
