# keepalived/logrotate-supervisor.ini
[program:logrotate]
command=/usr/local/bin/logrotate-script.sh
autostart=true
autorestart=true
